# src/gasnd.service
[Unit]
After=sshd.service
[Service]
ExecStart=/usr/local/sbin/gasnd
[Install]
WantedBy=multi-user.target
