# src/gdnsd.service
[Unit]
After=sshd.service
[Service]
ExecStart=/usr/local/sbin/gdnsd
[Install]
WantedBy=multi-user.target
