# src/ghistoryd.service
[Unit]
After=sshd.service
[Service]
ExecStart=/usr/local/sbin/ghistoryd
[Install]
WantedBy=multi-user.target
