# src/grdnsd.service
[Unit]
After=sshd.service
[Service]
ExecStart=/usr/local/sbin/grdnsd
[Install]
WantedBy=multi-user.target
